# systemd/.config/systemd/user/futon-backup.timer
[Unit]
Description=futon weekly backup timer

[Timer]
OnCalendar=Mon *-*-* 3:00
Persistent=true

[Install]
WantedBy=default.target
