# systemd/.config/systemd/user/home-software.automount
[Unit]
After=network-online.target
Before=remote-fs.target

[Automount]
Where=/home/software
TimeoutIdleSec=600

[Install]
WantedBy=default.target

