# systemd/.config/systemd/user/kanshi.service
[Unit]
Description=kanshi

[Service]
ExecStart=/usr/local/bin/kanshi
Restart=on-failure
