# systemd/.config/systemd/user/swayped.service
[Unit]
Description=swayped

[Service]
ExecStart=/home/adrien/.cargo/bin/swayped -S
Restart=on-failure
