# systemd/.config/systemd/user/udiskie.service
[Unit]
Description=udiskie automounter

[Service]
ExecStart=udiskie -T
Restart=on-failure

[Install]
WantedBy=default.target
